# systemd/mt5bridge-retention.timer
[Unit]
Description=Run MT5Bridge tick/candle retention daily at 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
