# systemd/mt5bridge.service
[Unit]
Description=MT5Bridge stack (Postgres, Redis, backend, Prometheus, Grafana via Docker Compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/mt5bridge
EnvironmentFile=/etc/mt5bridge/mt5bridge.env
ExecStart=/usr/bin/docker compose --env-file /etc/mt5bridge/mt5bridge.env up -d
ExecStop=/usr/bin/docker compose --env-file /etc/mt5bridge/mt5bridge.env down
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
