# systemd/mt5bridge-backup.timer
[Unit]
Description=Run MT5Bridge Postgres backup nightly at 02:00

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
